# sysconf/aqua-webapp.service
[Unit]
Description=Aquaponics Web Interface

[Service]
ExecStart=cd ~aqua/ap/webapp && ./start.R
User=aqua
Group=aqua
